# deploy-resources/practicalgowebsite.socket
[Socket]
ListenStream = 8080
BindIPv6Only = both

[Install]
WantedBy = sockets.target
